# contrib/init/bitcoind.conf
description "nealcoin Core Daemon"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

env nealcoinD_BIN="/usr/bin/nealcoind"
env nealcoinD_USER="nealcoin"
env nealcoinD_GROUP="nealcoin"
env nealcoinD_PIDDIR="/var/run/nealcoind"
# upstart can't handle variables constructed with other variables
env nealcoinD_PIDFILE="/var/run/nealcoind/nealcoind.pid"
env nealcoinD_CONFIGFILE="/etc/nealcoin/nealcoin.conf"
env nealcoinD_DATADIR="/var/lib/nealcoind"

expect fork

respawn
respawn limit 5 120
kill timeout 60

pre-start script
    # this will catch non-existent config files
    # nealcoind will check and exit with this very warning, but it can do so
    # long after forking, leaving upstart to think everything started fine.
    # since this is a commonly encountered case on install, just check and
    # warn here.
    if ! grep -qs '^rpcpassword=' "$nealcoinD_CONFIGFILE" ; then
        echo "ERROR: You must set a secure rpcpassword to run nealcoind."
        echo "The setting must appear in $nealcoinD_CONFIGFILE"
        echo
        echo "This password is security critical to securing wallets "
        echo "and must not be the same as the rpcuser setting."
        echo "You can generate a suitable random password using the following"
        echo "command from the shell:"
        echo
        echo "bash -c 'tr -dc a-zA-Z0-9 < /dev/urandom | head -c32 && echo'"
        echo
        echo "It is also recommended that you also set alertnotify so you are "
        echo "notified of problems:"
        echo
        echo "ie: alertnotify=echo %%s | mail -s \"nealcoin Alert\"" \
            "admin@foo.com"
        echo
        exit 1
    fi

    mkdir -p "$nealcoinD_PIDDIR"
    chmod 0755 "$nealcoinD_PIDDIR"
    chown $nealcoinD_USER:$nealcoinD_GROUP "$nealcoinD_PIDDIR"
    chown $nealcoinD_USER:$nealcoinD_GROUP "$nealcoinD_CONFIGFILE"
    chmod 0660 "$nealcoinD_CONFIGFILE"
end script

exec start-stop-daemon \
    --start \
    --pidfile "$nealcoinD_PIDFILE" \
    --chuid $nealcoinD_USER:$nealcoinD_GROUP \
    --exec "$nealcoinD_BIN" \
    -- \
    -pid="$nealcoinD_PIDFILE" \
    -conf="$nealcoinD_CONFIGFILE" \
    -datadir="$nealcoinD_DATADIR" \
    -disablewallet \
    -daemon

